# app/orin_nano/camera_service_new/camera-service.service
[Unit]
Description=Lightweight Buffer-Based Jetson Camera API Service
After=network.target

[Service]
Type=simple
User=azuolas
WorkingDirectory=/home/azuolas/mmoment/app/orin_nano/camera_service_new
# Run port cleanup script before starting to avoid conflicts
ExecStartPre=/bin/bash -c "lsof -ti:5003 | xargs -r kill -9"
ExecStart=/usr/bin/python3 /home/azuolas/mmoment/app/orin_nano/camera_service_new/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 