# app/orin_nano/systemd/frontend-bridge.service
[Unit]
Description=Jetson Camera Frontend Bridge
After=network.target
Wants=camera-service.service

[Service]
Type=simple
User=azuolas
WorkingDirectory=/home/azuolas/mmoment/app/orin_nano/frontend_bridge
ExecStart=/usr/bin/python3 /home/azuolas/mmoment/app/orin_nano/frontend_bridge/frontend_bridge.py --camera-port 5002 --frontend-port 5003
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 