# app/raspberry_pi/solana-middleware.service
[Unit]
Description=Solana Transaction Verification Middleware
After=network.target camera-api.service

[Service]
User=azuolas
Group=azuolas
WorkingDirectory=/home/azuolas/mmoment/app/raspberry_pi/solana-middleware
Environment="PATH=/home/azuolas/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/azuolas/mmoment/app/raspberry_pi/solana-middleware/venv/bin/python middleware.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 