# app/orin_nano/archive/systemd/cloudflared-compat.service
[Unit]
Description=Cloudflare Tunnel (compatible version)
After=network.target
Wants=frontend-bridge.service

[Service]
Type=simple
User=azuolas
WorkingDirectory=/home/azuolas
ExecStart=/usr/local/bin/cloudflared-compat tunnel run 6257e873-7943-4b85-b8a3-72b5b9d0a500
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target 